# service/hostpanel-nginx.service
[Unit]
Description=HostPanel Nginx
After=network.target

[Service]
Type=forking
PIDFile=/run/nginx.pid
ExecStart=/opt/hostpanel/plugins/nginx/nginx -p /opt/hostpanel/plugins/nginx -c /opt/hostpanel/plugins/nginx/nginx.conf
ExecReload=/opt/hostpanel/plugins/nginx/nginx -p /opt/hostpanel/plugins/nginx -s reload
ExecStop=/opt/hostpanel/plugins/nginx/nginx -p /opt/hostpanel/plugins/nginx -s stop
Restart=on-failure

[Install]
WantedBy=multi-user.target
